# templates/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/log/supervisor/supervisord.pid

[program:nginx]
command = nginx -g 'daemon off;'

[program:php-fpm]
command = php-fpm --nodaemonize
